# deploy/agent-service.service
[Unit]
Description=Agent Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/agent-service
EnvironmentFile=/home/ubuntu/agent-service/.env
ExecStart=/home/ubuntu/agent-service/venv/bin/python -m src.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
